# systemd/jafar-morning.timer
[Unit]
Description=Jafar Morning Run Timer (7am-12pm randomized)
Documentation=https://github.com/yourusername/twitter_sentiment_analysis
Requires=jafar.service

[Timer]
# Run every day at 7:00 AM
OnCalendar=*-*-* 07:00:00

# Add 0-5 hours of randomization (7am-12pm window)
RandomizedDelaySec=18000

# If the system was off when the timer should have triggered, run it on next boot
Persistent=true

# Ensure we don't run if evening run just completed
# (systemd will skip if service is already running or just ran)

[Install]
WantedBy=timers.target
